# 5-runit/sv/randomseed/run
#!/bin/bash -e

echo -e "---- ./run randomseed inicio"
exec 2>&1
exec /etc/sv/randomseed/randomseed-service_ini.sh

