# 5-runit/sv/template_4/run
#!/bin/bash -e

echo -e "---- ./run service inicio"
exec 2>&1
exec /etc/sv/network/network_ini.sh

