# 5-runit/sv/dbusdaemon/finish
#!/bin/bash -e

echo "Apagando dbus"
exec 2>&1

kill $(cat /run/dbus/pid)
rm -rf /run/dbus/system_bus_socket /run/dbus/pid
