# 5-runit/sv/localnet/finish
#!/bin/bash -e

echo -e "---- ./finish stopping localnet"
exec 2>&1
exec /etc/sv/localnet/localnet_fin.sh

