# 5-runit/sv/iptables/finish
#!/bin/bash -e

echo -e "iptables service finish"
exec 2>&1
exec /etc/sv/iptables/iptables_fin.sh

