# 5-runit/sv/iptables/run
#!/bin/bash -e

echo -e "---- ./run iptables inicio"
exec 2>&1
exec /etc/sv/iptables/iptables_ini.sh

