# 5-runit/sv/getty-3/run
#!/bin/sh 
echo "---- ./run Se abre la consola tty3"
exec /usr/local/bin/chpst -P /sbin/agetty 38400 tty3 linux
