# 5-runit/sv/syslogd/run
#!/bin/bash
echo "---- ./run syslogd!"
exec /sbin/syslogd -n -m 0

