# 5-runit/sv/network.wlan0/finish
#!/bin/bash -e

echo -e "---- ./finish network.wlan0"
exec 2>&1
exec /etc/sv/network.wlan0/network.wlan0_fin.sh

