# 5-runit/sv/network.wlan0/run
#!/bin/bash -e


/usr/local/bin/sv start iptables || exit 1

echo -e "---- ./run network wlan0 inicio"
exec 2>&1
exec /etc/sv/network.wlan0/network.wlan0_ini.sh

