# 5-runit/sv/udevd/run
#!/bin/bash -e
echo "daemon udevd"
exec 2>&1
exec /sbin/udevd 2>&1
