# 5-runit/sv/getty-1/run
#!/bin/sh
echo "---- ./run Se abre la consola tty1"
exec /usr/local/bin/chpst -P /sbin/agetty 38400 tty1 linux
