# 5-runit/sv/sshd/run
#!/bin/bash -e

echo -e "---- ./run sshd inicio"

exec 2>&1
/usr/sbin/sshd -D -e

echo -e "es despues de run sshd ¿me ves?"
