# 5-runit/sv/klogd/run
#!/bin/bash

echo "---- ./run klogd"
/usr/local/bin/sv up syslogd || exit 1
exec /sbin/klogd -n -c 4

